# Aria2_Ng/service/caddy2/run
#!/bin/sh

export HOME=/root XDG_CONFIG_HOME=/root

exec /caddy run --config /Caddyfile

