# Aria2_Ng/service/aria2/run
#!/bin/sh

umask 0

exec 1>/dev/null 2>&1

exec /aria2c --conf-path=/aria2.conf

